# scripts/odoo-backup.timer
[Unit]
Description=Run the Odoo (KEIT) backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
